# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/voice-stack/supervisord.log
pidfile=/tmp/supervisord.pid

[program:api]
directory=/opt/voice-stack
command=python3 -m uvicorn app.main:app --host 0.0.0.0 --port %(ENV_PORT)s
startsecs=3
autostart=true
autorestart=true
stdout_logfile=/var/log/voice-stack/api.log
stderr_logfile=/var/log/voice-stack/api.err.log

[program:stt]
directory=/opt/voice-stack
command=python3 -m uvicorn services.stt_nemotron_server:app --host 0.0.0.0 --port %(ENV_STT_PORT)s
autostart=true
autorestart=true
stdout_logfile=/var/log/voice-stack/stt.log
stderr_logfile=/var/log/voice-stack/stt.err.log

[program:brain]
command=/opt/voice-stack/scripts/start_brain.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/voice-stack/brain.log
stderr_logfile=/var/log/voice-stack/brain.err.log

[program:tts]
command=/opt/voice-stack/scripts/start_tts.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/voice-stack/tts.log
stderr_logfile=/var/log/voice-stack/tts.err.log
